# targetfs/common/50-pvrsrvctl.rules.template
SUBSYSTEM=="platform", ACTION=="bind", ENV{DRIVER}=="pvrsrvkm", RUN+="{{BINDIR}}/pvrsrvctl --start --no-module"
